# deployment/quoteflow.service
[Unit]
Description=QuoteFlow Backend Service
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/quoteflow/server
ExecStart=/usr/bin/node index.js
Restart=on-failure
EnvironmentFile=-/etc/quoteflow/quoteflow.env

[Install]
WantedBy=multi-user.target
